# DNS_Performance_MRTG/etc/cron.d/mrtg
*/5 * * * * root LANG=C LC_ALL=C /usr/bin/mrtg /etc/mrtg/dns.cfg --lock-file /var/lock/mrtg/mrtg_n --confcache-file /var/lib/mrtg/dns.ok
