# WAG200G_DSL_Conn_Status/etc/cron.d/mrtg
*/5 * * * * root LANG=C LC_ALL=C /usr/bin/mrtg /etc/mrtg/dsl.cfg --lock-file /var/lock/mrtg/mrtg_d --confcache-file /var/lib/mrtg/dsl.ok
